# automation/systemd/scheduled-message-replies.service.d/10-network-ready.conf
[Service]
ExecStartPre=/usr/bin/nm-online --wait-for-startup --quiet --timeout=60
ExecStartPre=/usr/bin/nm-online --quiet --timeout=60
